# rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::sysinit:-/bin/bash
